# etc/systemd/system/backup-sysinfo.timer
[Unit]
Description=Schedule system information backup.

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=basic.target
