# etc/systemd/system/backup-sysinfo.service
[Unit]
Description=Write system information to backup drive.
RequiresMountsFor=/backup

[Service]
ExecStart=/usr/local/bin/sysinfo.sh
